# deploy/supervisord.conf
; Process supervisor for the momentum-desk container. Three programs:
;   ib_proxy  — local HTTP->HTTPS re-proxy so the gateway's upstream TLS comes
;               from Python's stack (evades IBKR's JA3/JA4 WAF block). priority
;               100 = starts first.
;   ibeam     — headless Chromium auto-submits IBKR credentials into the CP
;               gateway login, then keeps the session alive. The only manual
;               step is approving the IBKR Mobile push (IB Key 2FA) on your phone.
;   dashboard — FastAPI/uvicorn on :8000, serves the API + the built Vite SPA.
;
; Run with `supervisord -n -c /etc/supervisor/supervisord.conf`.

[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ------------------------------------------------------------------
[program:ib_proxy]
command=/opt/venv/bin/python -m momentum_desk.ib_proxy
priority=100
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/app/logs/ib_proxy.out.log
stderr_logfile=/app/logs/ib_proxy.err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1",LOG_LEVEL="INFO"

; ------------------------------------------------------------------
; ibeam — IBKR_USERNAME / IBKR_PASSWORD / IBKR_PAPER come from Fly secrets
; (the same names used everywhere); the lines below map them into the IBEAM_*
; names ibeam reads. IBEAM_MAX_FAILED_AUTH=2 is a circuit breaker: after two
; failed logins ibeam stops rather than burning retries and locking the account.
[program:ibeam]
command=/opt/venv/bin/python /srv/ibeam/ibeam_starter.py
priority=200
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
stdout_logfile=/app/logs/ibeam.out.log
stderr_logfile=/app/logs/ibeam.err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
environment=
    PYTHONUNBUFFERED="1",
    IBEAM_ACCOUNT="%(ENV_IBKR_USERNAME)s",
    IBEAM_PASSWORD="%(ENV_IBKR_PASSWORD)s",
    IBEAM_SECRETS_SOURCE="env",
    IBEAM_GATEWAY_DIR="/opt/gateway",
    IBEAM_CHROME_DRIVER_PATH="/usr/bin/chromedriver",
    IBEAM_GATEWAY_BASE_URL="https://localhost:5000",
    IBEAM_USE_PAPER_ACCOUNT="%(ENV_IBKR_PAPER)s",
    IBEAM_MAX_FAILED_AUTH="2",
    IBEAM_MAX_IMMEDIATE_ATTEMPTS="3",
    IBEAM_OAUTH_TIMEOUT="300",
    IBEAM_PAGE_LOAD_TIMEOUT="60",
    IBEAM_AUTHENTICATION_STRATEGY="B",
    IBEAM_MAX_STATUS_CHECK_RETRIES="300",
    IBEAM_MAX_REAUTHENTICATE_RETRIES="6",
    IBEAM_IBKEY_PROMO_EL_CLASS="CLASS_NAME@@ibkey-promo-skip"

; ------------------------------------------------------------------
; Dashboard — bound to 0.0.0.0 so Fly's edge proxy reaches it; basic auth
; (DASHBOARD_USER/PASSWORD) guards every route except /api/health. IBKR_ENABLED
; turns on the in-process keepalive + /api/ibkr/status (see server.py lifespan).
[program:dashboard]
command=/opt/venv/bin/uvicorn momentum_desk.server:app --host 0.0.0.0 --port 8000
priority=300
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/app/logs/dashboard.out.log
stderr_logfile=/app/logs/dashboard.err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1",IBKR_ENABLED="true"

; ------------------------------------------------------------------
; Autotrade — the autonomous paper loop. OFF BY DEFAULT (autostart=false): a
; deploy never starts trading. Turn it on deliberately once you trust it:
;   fly ssh console --app momentum-desk-mav
;   supervisorctl start autotrade
; Requires DATA_FEED=polygon (real candidates) and an authenticated paper
; gateway. --send actually transmits (paper); drop it for a dry rehearsal.
; Drop a /app/data/HALT file (the mounted volume) to flatten + stop.
[program:autotrade]
command=/opt/venv/bin/python -m scripts.live_paper --broker cp --send --poll 30
priority=400
autostart=false
autorestart=true
startsecs=10
stopwaitsecs=30
stdout_logfile=/app/logs/autotrade.out.log
stderr_logfile=/app/logs/autotrade.err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1",HALT_FILE="/app/data/HALT"
